# files/fpm.conf
[program:php]
command = /usr/sbin/php5-fpm --nodaemonize
user = root
